# supervisor/warp.ini
[program:warp-worker]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "/bin/warp -l ${LISTEN_ADDRESS} -o ${LISTEN_PORT} --listen-protocol ${LISTEN_PROTOCOL} --remote-address ${REMOTE_ADDRESS} --remote-port ${REMOTE_PORT} --remote-protocol ${REMOTE_PROTOCOL} --threads ${MAX_THREADS} --timeout ${TIMEOUT}"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/warp.log
stopwaitsecs=3600